# supervisor/kafka_manager.conf
[program:kafka-manager]
command=python /opt/kafka_manager/KafkaManager.py
autostart=true
autorestart=true