# supervisor/kafka_manager.conf
[program:kafka-manager]
command=python %(KAFKA_MANAGER_HOME)s/KafkaManager.py
autostart=true
autorestart=true